# ios/layout/Library/LaunchDaemons/dev.r58playz.kdeconnectjb.plist
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.r58playz.kdeconnectjb</string>
    <key>Program</key>
    <string>/var/jb/usr/local/bin/kdeconnectd</string>

    <key>UserName</key>
    <string>mobile</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
